# bootc-ucore/etc/systemd/system/bootc-fetch-apply-updates.timer.d/weekly.conf
[Timer]
RandomizedDelaySec=600
#Apply updates on Saturday at 4am
OnCalendar=Fri *-*-* 04:00:00

